# fpkg_criticalupdated/criticalupdated
#!/sbin/openrc-run
extra_commands=""
depend(){
	need net
	provide updctl
}
command="@bindir@/fpkg"
command_args="updated"
pidfile="/run/${RC_SVCNAME}.pid"
